# Projects/tools/LinuxHost/configuration/BeagleBoneBlack/upstart/remoti_uinput.conf
start on stopped remoti_init
stop on stopping npi_server

# This application is no daemon, so no need for expect stanza
# expect daemon

console log

exec /etc/init.d/remoti_uinput -D 127.0.0.1:2533
