# Projects/tools/LinuxHost/configuration/systemd/rf4ce.service.in
[Unit]
Description=@WITBE_SYSTEMD_SERVICE_NAME@
After=@WITBE_SYSTEMD_SERVICE_RUN_AFTER@

[Service]
Type=@WITBE_SYSTEMD_SERVICE_TYPE@
ExecStart=@WITBE_SYSTEMD_SERVICE_COMMAND@ @WITBE_SYSTEMD_SERVICE_ARGS@
Restart=@WITBE_SYSTEMD_SERVICE_RESTART@

[Install]
WantedBy=multi-user.target
